# dotfiles/.config/systemd/user/mbsync.service
[Unit]
Description=mbsync all accounts

[Service]
Type=oneshot
ExecStart=/home/nacef/.local/bin/mail_notify
