# dotfiles/.config/systemd/user/mbsync.service
[Unit]
Description=mbsync all accounts

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -a
